# org.eclipse.neoscada.contrib.assembly/kafka-bridge/deb/control/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable scada.app.kafka-bridge
systemctl start scada.app.kafka-bridge

exit 0