# org.eclipse.neoscada.contrib.kafka/deb/control/prerm
#!/bin/sh

systemctl stop scada.app.kafka-bridge
systemctl disable scada.app.kafka-bridge

exit 0