# org.eclipse.neoscada.contrib.iec104torest/deb/etc/init/scada.app.iec104torest.conf
# start iec104torest service

author "Juergen Rose"
description "service which converts exports a 104 connection as a rest endpoint"

start on (filesystem and runlevel [2345])
stop on runlevel [016]

respawn

expect daemon

env ARGS="-Dlogback.configurationFile=/var/lib/eclipsescada/iec104torest/logback.xml -jar /var/lib/eclipsescada/iec104torest/iec104torest.jar /var/lib/eclipsescada/iec104torest/config.json > /var/log/eclipsescada/scada.app.iec104torest-stdout.log 2> /var/log/eclipsescada/scada.app.iec104torest-stderr.log"

script
    exec start-stop-daemon --start -c eclipsescada -g eclipsescada --make-pidfile --pidfile /var/run/scada.app.iec104torest.pid --background --exec /usr/bin/java -- $ARGS
end script

pre-stop script
    exec start-stop-daemon --signal QUIT --stop -c eclipsescada -g eclipsescada --make-pidfile --pidfile /var/run/scada.app.iec104torest.pid --exec /usr/bin/java -- $ARGS
end script
